# examples/3.X/root/etc/services.d/app/run
#!/usr/bin/with-contenv sh

exec python /app/example.py
