# examples/2.7/root/etc/services.d/app/run
#!/usr/bin/with-contenv sh

exec python2 /app/example.py
